# docker-complete-solution/ml_docker/supervisor.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:ml-detection]
command=/usr/src/app/ml.sh
autorestart=true
